# infra/ansible/roles/docker_host/templates/tea-party-compose.service.j2
[Unit]
Description=Tea Party Reservation Bot Compose Stack ({{ item.key }})
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ env_cfg.compose_dir }}
EnvironmentFile=-{{ env_cfg.env_dir }}/app.env
ExecStart=/usr/bin/docker compose -f {{ env_cfg.compose_dir }}/compose.yml -f {{ env_cfg.compose_dir }}/{{ env_cfg.compose_override }} up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f {{ env_cfg.compose_dir }}/compose.yml -f {{ env_cfg.compose_dir }}/{{ env_cfg.compose_override }} down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
